# scripts/freezerbot.service
[Unit]
Description=Freezerbot Voice Inventory System
After=network-online.target sound.target
Wants=network-online.target
ConditionPathExists=/media/pi/FREEZERBOT/bootstrap.sh

[Service]
User=pi
Group=pi
Environment=FREEZERBOT_ROOT=/media/pi/FREEZERBOT
Environment=QT_QPA_PLATFORM=xcb
WorkingDirectory=/media/pi/FREEZERBOT
ExecStart=xinit /bin/bash /media/pi/FREEZERBOT/bootstrap.sh -- :0 vt1 -nocursor
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
